# python/Sytem_specification/CPE-master (3)/CPE-master/SystemSetup/telegraf.service
[Unit]
Description=Telegraf service

[Service]
ExecStart=/usr/bin/telegraf -config /etc/telegraf/telegraf.conf
NotifyAccess=main
#WatchdogSec=10
Restart=on-failure
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target